# script/mysql/supervisord-mysqld.conf
[program:mysqld]
command=/script/mysql/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true